# config/launchd/com.user.bugwarrior-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.bugwarrior-pull</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>bugwarrior</string>
        <string>pull</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/emiller/src/personal/taskagent/bugwarrior</string>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/emiller/.local/share/bugwarrior/pull.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/emiller/.local/share/bugwarrior/pull-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/emiller/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>