# config/launchd/com.user.taskwarrior-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.taskwarrior-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/task</string>
        <string>sync</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/emiller/.local/share/taskwarrior/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/emiller/.local/share/taskwarrior/sync-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>